# roles/awesomewm/files/systemd/ai-usage-monitor.service
[Unit]
Description=AI Usage Monitor - polls the active AI provider

[Service]
Type=simple
ExecStart=%h/.local/bin/ai-usage-monitor.sh
Restart=on-failure
RestartSec=10
KillSignal=SIGTERM
TimeoutStopSec=5

[Install]
WantedBy=default.target
